# etc/xdg/autostart/teqball.desktop
[Desktop Entry]
Type=Application
Name=TeqBall
Exec=sudo /usr/bin/python3 /home/pi/TeqBall/main.py
